# deployharness.service
[Unit]
Description=AI Harness — personal AI studio
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/hAIrness
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
